# rfid_dash_dev.service
[Unit]
Description=RFID Dashboard Flask App
After=network.target

[Service]
User=tim
Group=www-data
WorkingDirectory=/home/tim/test_rfidpi
Environment="PATH=/home/tim/test_rfidpi/venv/bin"
ExecStart=/home/tim/test_rfidpi/venv/bin/gunicorn --workers 1 --bind 0.0.0.0:7409 --chdir /home/tim/test_rfidpi run:app
Restart=always

[Install]
WantedBy=multi-user.target