# rfid_dash_dev.service
[Unit]
Description=RFID Dashboard Flask App
After=network.target mariadb.service

[Service]
User=tim
Group=www-data
Environment=PROJECT_DIR=/opt/RFID3
WorkingDirectory=${PROJECT_DIR}
Environment="PATH=${PROJECT_DIR}/venv/bin:/usr/bin"
Environment="APP_IP=192.168.3.110"
ExecStart=${PROJECT_DIR}/venv/bin/gunicorn --workers 1 --threads 4 --timeout 600 --bind 0.0.0.0:8102 --error-logfile ${PROJECT_DIR}/logs/gunicorn_error.log --access-logfile ${PROJECT_DIR}/logs/gunicorn_access.log run:app
ExecStop=/bin/kill -s KILL $MAINPID
Restart=always
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

