# scripts/laundry-auto-finalize.service
[Unit]
Description=RFID3 Laundry Auto-Finalization Service
Wants=laundry-auto-finalize.timer

[Service]
Type=oneshot
User=tim
WorkingDirectory=/home/tim/RFID3
Environment=PATH=/home/tim/RFID3/venv/bin:/usr/bin
ExecStart=/home/tim/RFID3/venv/bin/python /home/tim/RFID3/scripts/auto_finalize_laundry.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target