# scripts/rfid-auto-update.service
[Unit]
Description=RFID3 Auto Update Service
Wants=rfid-auto-update.timer

[Service]
Type=oneshot
User=root
Environment=PROJECT_DIR=/opt/RFID3
ExecStart=${PROJECT_DIR}/scripts/auto_update.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target